# famous.crontab
0 11 * * * cd PATH_TO_YOUR_FOLDER && docker compose up
